# packaging/deb/etc/systed/system/forecastmetrics.service
[Unit]
Description = ForecastMetrics

[Service]
Type = simple
User = pi
ExecStart = /usr/local/bin/forecastmetrics
