# packaging/deb/etc/systed/system/forecastmetrics.timer
[Unit]
Description = ForecastMetrics Timer

[Timer]
OnCalendar = *-*-* *:00:00 America/New_York

[Install]
WantedBy = timers.target